# packaging/systemd/tires.service
[Unit]
Description=Tires - Tiered Storage Manager
Documentation=https://github.com/tires/tires/blob/main/README.md
After=local-fs.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tires /etc/tires/storage.json
WorkingDirectory=/var/lib/tires
StandardOutput=journal
StandardError=journal
User=root
Group=root

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/tires
PrivateTmp=true

[Install]
WantedBy=multi-user.target
