# packaging/systemd/tires.timer
[Unit]
Description=Run Tires - Tiered Storage Manager
Documentation=https://github.com/gailoks/tires/blob/main/README.md

[Timer]
# Run interval is configured in /etc/tires/storage.json (RunInterval field)
# Supported values: "minutely", "hourly", "daily", "weekly", "monthly"
# or systemd calendar format (e.g., "*-*-* 02:00:00" for daily at 2 AM)
# Default: hourly
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
